# lastseen.crontab
*/2	*	*	*	*	/prj/yuridm.github.io/lastseen.sh >> /prj/yuridm.github.io/lastseen.log 2>&1
